# launchd/com.auto-report.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auto-report.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/xiaozhangxuezhang/Documents/GitHub/auto-daily-report/scripts/run-claude-tmux.sh</string>
        <string>weekly-report</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/xiaozhangxuezhang/Documents/GitHub/auto-daily-report/logs/launchd-weekly.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/xiaozhangxuezhang/Documents/GitHub/auto-daily-report/logs/launchd-weekly.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/xiaozhangxuezhang/.npm-global/bin</string>
        <key>HOME</key>
        <string>/Users/xiaozhangxuezhang</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/xiaozhangxuezhang/Documents/GitHub/auto-daily-report</string>
</dict>
</plist>
